# packaging/debian/ledger-watcher.path
[Unit]
Description=openbank transaction ledger configuration watcher
After=ledger.service

[Path]
PathChanged=/etc/ledger/conf.d
Unit=ledger-watcher.service

[Install]
WantedBy=paths.target
