# packaging/debian_armhf/lib/systemd/system/ledger-unit@.service
[Unit]
Description=openbank transaction ledger tenant %i
PartOf=ledger.service
After=ledger.service
ConditionPathExists=/etc/init/ledger.conf

[Service]
Type=notify
Environment=LEDGER_TENANT=%i
EnvironmentFile=/etc/init/ledger.conf
StartLimitBurst=0
ExecStart=/usr/bin/ledger-unit %i
StandardInput=null
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
Restart=always
RestartSec=0
TimeoutStartSec=0
TimeoutSec=infinity
KillMode=process

[Install]
WantedBy=ledger.service
