# packaging/debian_arm64/lib/systemd/system/ledger.path
[Unit]
Description=openbank transaction ledger services namespace
PartOf=ledger.service
ConditionPathExists=/etc/init/ledger.conf

[Path]
PathChanged=/etc/init/ledger.conf
PathModified=/etc/init/ledger.conf
Unit=ledger.service

[Install]
WantedBy=ledger.service
